# xrayvpn/debian/xrayvpnd-refresh.service
[Unit]
Description=Refresh xrayvpn geodata and RU CIDRs
After=xrayvpnd.service
Requires=xrayvpnd.service

[Service]
Type=oneshot
ExecStart=/usr/bin/xrayvpn refresh
User=xrayvpn

