# templates/hstd/clientrotate.service.j2
[Unit]
Description=Rotate xray client UUID

[Service]
Type=oneshot
ExecStart=/usr/bin/xrayvpn rotate {{ host.data.rotate_secret }}
