# templates/hstd/clientrotate.timer.j2
[Unit]
Description=Daily xray client UUID rotation

[Timer]
OnCalendar=*-*-* 03:00,06,12,18,24,30,36,42,48,54:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
