# templates/xray_server/clientrotate.service.j2
[Unit]
Description=Rotate xray client UUID

[Service]
Type=oneshot
ExecStart=/usr/local/bin/clientrotate {{ host.data.rotate_secret }}
ExecStartPost=/bin/systemctl restart xray
